# resources/deb-scripts/prerm
#!/bin/bash
/bin/systemctl stop bungee-backup
/bin/systemctl disable bungee-backup